# logstash-shipper.upstart
description "logstash shipper server"
 
start on runlevel [2345]
stop on runlevel [06]
 
# tell upstart we're creating a daemon
# upstart manages PID creation for you.
expect fork
 
script
  exec sudo /opt/logstash/bin/logstash agent -f /etc/logstash-shipper.conf &
  emit logstash-shipper-server_running
end script